# etc/systemd/system/http.service
[Unit]
Description=http

[Service]
ExecStart=/usr/bin/http

[Install]
WantedBy=multi-user.target
